# systemd/solarinspector-updater.path
[Unit]
Description=Watch for SolarInspector update requests

[Path]
PathExists=/var/lib/solarinspector/update-request.json
Unit=solarinspector-updater.service

[Install]
WantedBy=multi-user.target
