# systemd/solarinspector-updater.service
[Unit]
Description=SolarInspector Privileged Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root

Environment=PYTHONPATH=/opt/solarinspector/current/app

ExecStart=/opt/solarinspector/current/.venv/bin/python \
  /opt/solarinspector/current/app/updater_service.py

PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict

ReadWritePaths=/opt/solarinspector
ReadWritePaths=/var/lib/solarinspector
ReadWritePaths=/var/cache/solarinspector
ReadWritePaths=/var/log/solarinspector

[Install]
WantedBy=multi-user.target
